# services/executer-api/usbswitcher.service
[Unit]
Description=Offers Api to toggle the usb port power
After=network.target
StartLimitIntervalSec=5

[Service]
Type=simple
Restart=always
RestartSec=10
Group=pi
User=pi
Environment="SER_DIR=/home/pi/my-remote-usb-controller/services/executer-api"
ExecStart=source "${SER_DIR}/venv/bin/activate" && python "${SER_DIR}/src/main.py"

[Install]
WantedBy=multi-user.target