# services/executer-api/usbswitcher.service
[Unit]
Description=Offers Api to toggle the usb port power
After=network.target
StartLimitIntervalSec=5

[Service]
Type=simple
Restart=always
RestartSec=10
Group=pi
User=pi
ExecStart=./run.sh

[Install]
WantedBy=multi-user.target